# pyspark-docker/supervisor/pyspark.conf
[program:pyspark]
command=/usr/bin/start-pyspark.sh
autostart=true
autorestart=true